# config/supervisord/desktop.conf
[supervisord]
nodaemon=true

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman -nodaemon
autostart=true
autorestart=true
priority=100

[program:xrdp]
command=/usr/sbin/xrdp -nodaemon
autostart=true
autorestart=true
priority=200

[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true